# tools/cron-vps/systemd/keymatch-cron-visites-rappel.timer
[Unit]
Description=KeyMatch cron timer visites-rappel — schedule "0 10,18 * * *"

[Timer]
OnCalendar=*-*-* 10,18:00:00
Persistent=true
RandomizedDelaySec=2min
AccuracySec=1min

[Install]
WantedBy=timers.target
