# tools/cron-vps/systemd/keymatch-cron-post-bail.timer
[Unit]
Description=KeyMatch cron timer post-bail — schedule "0 9 * * *"

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
RandomizedDelaySec=2min
AccuracySec=1min

[Install]
WantedBy=timers.target
