# tools/cron-vps/systemd/keymatch-cron-import-jobs-cleanup.service
[Unit]
Description=KeyMatch cron import-jobs-cleanup — Cleanup jobs import abandonnés
After=docker.service network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/keymatch/NextImmo-main
ExecStart=/bin/bash /opt/keymatch/NextImmo-main/tools/cron-vps/scripts/run-cron.sh import-jobs-cleanup /api/cron/import-jobs-cleanup
StandardOutput=append:/var/log/keymatch-cron.log
StandardError=append:/var/log/keymatch-cron.log
# Si fail, ne pas spammer (cron-routes.tsv est source de vérité)
Restart=no
TimeoutStartSec=320s

[Install]
WantedBy=multi-user.target
