# tools/cron-vps/systemd/keymatch-cron-fetcher-health.timer
[Unit]
Description=KeyMatch cron timer fetcher-health — schedule "0 */6 * * *"

[Timer]
OnCalendar=*-*-* 00/06:00:00
Persistent=true
RandomizedDelaySec=2min
AccuracySec=1min

[Install]
WantedBy=timers.target
