# tools/cron-vps/systemd/keymatch-cron-edl-contestation-retard.timer
[Unit]
Description=KeyMatch cron timer edl-contestation-retard — schedule "45 8 * * *"

[Timer]
OnCalendar=*-*-* 08:45:00
Persistent=true
RandomizedDelaySec=2min
AccuracySec=1min

[Install]
WantedBy=timers.target
