# tools/cron-vps/systemd/keymatch-cron-loyers-mensuels-create.timer
[Unit]
Description=KeyMatch cron timer loyers-mensuels-create — schedule "0 6 1 * *"

[Timer]
OnCalendar=*-*-01 06:00:00
Persistent=true
RandomizedDelaySec=2min
AccuracySec=1min

[Install]
WantedBy=timers.target
