# tools/cron-vps/systemd/keymatch-cron-candidatures-digest.timer
[Unit]
Description=KeyMatch cron timer candidatures-digest — schedule "30 19 * * *"

[Timer]
OnCalendar=*-*-* 19:30:00
Persistent=true
RandomizedDelaySec=2min
AccuracySec=1min

[Install]
WantedBy=timers.target
